# sys-apps/apparmor-systemd/files/apparmor.target
[Unit]
Description=AppArmor target
DefaultDependencies=no
Before=sysinit.target

[Install]
WantedBy=sysinit.target
